# installation/bin/webserver.service
[Unit]
Description=Servidor Web
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/npm run start --prefix /opt/SmarTAP/webguiSONDAS/webgui/
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/bin:/opt/SmarTAP/bin

[Install]
WantedBy=multi-user.target

